# packages/src/moodle-provas/DEBIAN/postrm
#!/bin/bash
#set -x

update-rc.d -f moodle-provas remove >/dev/null
update-rc.d -f iptables-persistent remove >/dev/null
